# c/apps/network/echod/etc/systemd/echod.service.in
[Unit]
Description=Echo daemon
Documentation=man:echod(8) man:echod.conf(5)
Wants=network.target
After=network.target

[Service]
Type=forking
ExecStart=/usr/sbin/echod
ExecReload=/bin/kill -HUP $MAINPID
Restart=on-failure

[Install]
WantedBy=multi-user.target
