# c/apps/network/echod/etc/echod.service.in
[Unit]
Description=Echo daemon
After=munge.service network.target remote-fs.target
ConditionPathExists=@sysconfdir@/echod.conf

[Service]
Type=forking
EnvironmentFile=-/etc/sysconfig/echod
ExecStart=@sbindir@/echod $ECHOD_OPTIONS
ExecReload=/bin/kill -HUP $MAINPID
PIDFile=/var/run/echod.pid
KillMode=process
LimitNOFILE=131072
LimitMEMLOCK=infinity
LimitSTACK=infinity
Delegate=yes
@SYSTEMD_TASKSMAX_OPTION@

[Install]
WantedBy=multi-user.target