# python/flask/examples/minitwit/supervisord.conf
[supervisord]
logfile=supervisord.log      ; (main log file;default $CWD/supervisord.log)
logfile_maxbytes=50MB        ; (max main logfile bytes b4 rotation;default 50MB)
logfile_backups=10           ; (num of main logfile rotation backups;default 10)
loglevel=info                ; (log level;default info; others: debug,warn,trace)
pidfile=supervisord.pid      ; (supervisord pidfile;default supervisord.pid)
nodaemon=false               ; (start in foreground if true;default false)
minfds=1024                  ; (min. avail startup file descriptors;default 1024)
minprocs=200                 ; (min. avail process descriptors;default 200)

[program:gunicorn]
command=/home/tfarina/src/repos/github.com/tfarina/rsc/python/flask/examples/minitwit/start_server.sh
directory=/home/tfarina/src/repos/github.com/tfarina/rsc/python/flask/examples/minitwit
user=tfarina
redirect_stderr=true
stdout_logfile=/home/tfarina/supervisord_stdout.log
stdout_logfile_maxbytes=20MB
stdout_logfile_backups=10
